# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
pidfile=/var/run/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:wikicontroller]
command=/usr/bin/node /controller/index.js
autostart=true
autorestart=true
stopasgroup=true
stderr_logfile=/var/log/wikicontroller.err.log
stdout_logfile=/var/log/wikicontroller.out.log

[program:livequery]
command=/usr/bin/node /controller/parse-live-query.js
autostart=true
autorestart=true
stopasgroup=true
stderr_logfile=/var/log/livequery.err.log
stdout_logfile=/var/log/livequery.out.log